# crontab
* * * * * sh /var/www/html/cron/vtigercron.sh
